# util/docker/web/runit/beanstalkd/run
#!/bin/bash

sudo apt update && sudo apt install beanstalkd
echo 'Spinning up Beanstalkd process...'

setuser azuracast beanstalkd -p 11300
